# supervisord.conf
[supervisord]
nodaemon=true

[program:telegram]
command=python app.py
stdout_logfile=/telegram.log